# support/dustyis.conf
description "Dusty.Is server"

start on startup
start on filesystem and started networking
stop on shutdown
respawn

chdir /home/dusty/dusty.is 
env NODE_ENV=production
env PORT=3000
exec npm start 
